# backend/backend.cron

0 12 * * * docker exec videosdb_backend-sadhguru_1 nice -n 20 poetry run python manage.py videosdb --check-for-new-videos

